# i3config/.config/systemd/user/lspmux.service
[Unit]
Description=Rust Analyzer (LSP Mux)

[Service]
ExecStart=/usr/bin/lspmux server
Restart=always

[Install]
WantedBy=default.target
